# build.assets/macos/tsh/tsh.app/Contents/Library/LaunchDaemons/com.gravitational.teleport.tsh.vnetd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.gravitational.teleport.tsh.vnetd</string>
	<key>BundleProgram</key>
	<string>Contents/MacOS/tsh</string>
	<key>ProgramArguments</key>
	<array>
		<string>tsh</string>
		<string>vnet-daemon</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.gravitational.teleport.tsh.vnetd</key>
		<true/>
	</dict>
	<!-- VNet daemon uses os_log rather than logging to stdout or stderr, but these keys are kept in
	case something does gets written there before the os_log logger is set up. -->
	<key>StandardErrorPath</key>
	<string>/var/log/vnet.log</string>
	<key>StandardOutPath</key>
	<string>/var/log/vnet.log</string>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>EnvironmentVariables</key>
	<dict>
		<!-- Auto-updates are disabled because $HOME isn't available,
		but let's also disable them explicitly anyway. -->
		<key>TELEPORT_TOOLS_VERSION</key>
		<string>off</string>
	</dict>
</dict>
</plist>
