# assets/aws/files/system/teleport-get-cert.service
[Unit]
Description=Service getting teleport certificates
ConditionPathExists=/etc/teleport.d/use-letsencrypt

[Service]
Type=oneshot
ExecStartPre=/usr/local/bin/teleport-lock
ExecStart=/usr/local/bin/teleport-get-cert