# assets/marketplace/single-ami/files/system/teleport-renew-cert.timer
[Unit]
Description=Timer rotating teleport certificates
ConditionPathExists=/etc/teleport.d/use-letsencrypt

[Timer]
OnBootSec=20min
OnCalendar=hourly
Persistent=true