# assets/marketplace/cloudformation/files/system/teleport-get-cert.service
[Unit]
Description=Service getting teleport certificates

[Service]
Type=oneshot
ExecStartPre=/usr/bin/teleport-lock
ExecStart=/usr/bin/teleport-get-cert
