# assets/systemd-unit-upgrader/teleport-upgrade.service
[Unit]
Description=sdt unit launched by timer

[Service]
ExecStart=/usr/local/bin/teleport-upgrade
