# assets/systemd-unit-upgrader/teleport-upgrade.service
[Unit]
Description=automatic teleport upgrader

[Service]
ExecStart=/usr/local/bin/teleport-upgrade
